# nginx/etc/systemd/system/generate_ports_json.timer
[Unit]
Description=Run generate_ports_json every 15 seconds

[Timer]
OnBootSec=15s
OnUnitActiveSec=15s
Unit=generate_ports_json.service

[Install]
WantedBy=timers.target
